# codeur-watcher.service
# Service systemd pour faire tourner Codeur Watcher en continu (VPS OVH).
# À copier dans /etc/systemd/system/codeur-watcher.service
# (adapté OVH : utilisateur "administrator", dossier /home/administrator/codeur-watcher)

[Unit]
Description=Codeur Watcher - veille codeur.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/home/administrator/codeur-watcher
ExecStart=/usr/bin/node index.js --watch
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
